# etc/crontabs/root
* * * * * root cd /var/www/drupal7 && drush cron
